# crontab
* * * * * /bin/bash -c "/opt/check-temp.sh &> /var/log/cron.log"
